# systemd/ptt-whisper.service
[Unit]
Description=PTT Whisper Voice Input
After=graphical-session.target whisper-server.service
Wants=whisper-server.service

[Service]
Type=simple
Environment="DISPLAY={{DISPLAY}}"
Environment="XAUTHORITY={{XAUTHORITY}}"
ExecStart={{PROJECT_DIR}}/venv/bin/python {{PROJECT_DIR}}/voice_input.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
