# systemd/whisper-server.service
[Unit]
Description=Whisper Speech Recognition Server
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}
ExecStart={{PROJECT_DIR}}/venv/bin/uvicorn whisper_server:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
